# pkg/deb/postinst
#!/bin/sh
set -e

# Post-installation script for space-captain-server

case "$1" in
    configure)
        # Create data directory if it doesn't exist
        if [ ! -d /var/lib/space-captain ]; then
            mkdir -p /var/lib/space-captain
            chown nobody:nogroup /var/lib/space-captain
            chmod 755 /var/lib/space-captain
        fi
        
        # Enable and start the systemd service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable space-captain-server.service >/dev/null 2>&1 || true
            systemctl start space-captain-server.service >/dev/null 2>&1 || true
            echo "Space Captain Server service enabled and started"
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0